# deploy/napoleon.service
[Unit]
description=napoleon gameserver

[Service]
ExecStart=/etc/napoleon/bin

[Install]
WantedBy=multi-user.target
